# ansible/roles/docker-restart-on-boot/templates/docker-restart.timer.j2
[Unit]
Description=Run Docker restart 10 minutes after boot

[Timer]
OnBootSec=10min
Unit=docker-restart.service

[Install]
WantedBy=timers.target
